# hotkey-handler/hotkey-handler.service
[Unit]
Description=Hotkey Handler (Camera + Copilot)
After=multi-user.target

[Service]
ExecStart=/usr/local/bin/hotkey-handler.py
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
