# system/etc/systemd/system/auto-brightness.service
[Unit]
Description=Keyboard backlight inverse to screen brightness
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/auto-brightness.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
