# system/etc/systemd/system/kbd-backlight-auto.service
[Unit]
Description=Auto keyboard backlight (inverse of screen brightness)
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/local/bin/kbd-backlight-auto.sh
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
